# init/repo_server.conf
description "repo_server"

start on runlevel [345]
stop on runlevel [!345]

console log

respawn

setuid repo
setgid repo
chdir /var/lib/repo

exec repo_server
