# package/template.upstart
description "A Dynamic Ubuntu Repository Server"

start on runlevel [2345]
stop on runlevel [!2345]

console log

setuid @@user@@
chdir @@datadir@@

respawn

exec "/usr/bin/@@name@@" -dir="@@datadir@@"
