# package/template.prerm
#!/bin/bash

# Make sure dpkg gets notified of any errors
set -e

if [ -e "/lib/systemd/system/@@package@@.service" ]; then
    systemctl stop @@package@@.service >/dev/null || true
    systemctl disable @@package@@.service >/dev/null || true
fi

# Success - reached the end
exit 0
