# package/template.postinst
#!/bin/bash

# Make sure dpkg gets notified of any errors
set -e

chown -R @@user@@:@@group@@ "@@datadir@@"

if [ -e "/lib/systemd/system/@@package@@.service" ]; then
    # Start the service running
    systemctl --system daemon-reload >/dev/null || true
    systemctl enable @@package@@.service >/dev/null
    systemctl start @@package@@.service >/dev/null
fi

# Success - reached the end
exit 0
